# key-stroke-counter.service
[Unit]
Description=key-stroke-counter

[Service]
ExecStart=/path/to/key-stroke-counter/binary -f /root/keys.log

[Install]
WantedBy=multi-user.target